# deploy/systemd/sidai-backend.service
[Unit]
Description=SID.AI Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/sidai/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/sidai/.env
ExecStart=/opt/sidai/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8888 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
